# lite/ipk/overlay/etc/systemd/system/advantagescope-lite.service
[Unit]
Description=AdvantageScope Lite Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/local/bin/advantagescope-lite
ExecStart=/usr/bin/env python3 lite_server.py
Restart=always
RestartSec=5
TimeoutStartSec=3
TimeoutStopSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target